# crontab
#!/bin/sh

cd /srv && php artisan schedule:run
